# taxibot.service
[Unit]
Description=Taxi Bot - Telegram Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/taxibot
Environment="PATH=/opt/taxibot/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/taxibot/venv/bin/python -m bot.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
